# etc/supervisor/conf.d/sync.conf
[program:quaran-sync]
command=bash -c '/dockerstartup/sync.sh'

autorestart=true
directory=/root
user=www-data