# etc/supervisor/apashe.conf
[program:apache20]
command=/usr/sbin/apache2ctl -DFOREGROUND