# etc/re.conf
[program:re-sync]
command=bash -c '/dockerstartup/rerun.sh'

directory=/dockerstartup/

user=root







